# deployment/cloudflare-tunnel.service
[Unit]
Description=Cloudflare Quick Tunnel for Ladoo Metrics
After=network.target nginx.service

[Service]
Type=simple
User=root
# Quick Tunnel: no Cloudflare account or domain needed.
# Gives a random public https://<random>.trycloudflare.com URL.
# The assigned URL is printed to the journal on startup — check with:
#   journalctl -u cloudflare-tunnel -n 20 | grep trycloudflare
ExecStart=/usr/local/bin/cloudflared tunnel --url http://localhost:80
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloudflared

[Install]
WantedBy=multi-user.target
